# systemd/kovo.service
[Unit]
Description=KOVO AI Agent
After=network.target redis.service
Wants=redis.service
[Service]
Type=simple
User=esam
WorkingDirectory=/opt/kovo
EnvironmentFile=/opt/kovo/config/.env
ExecStart=/opt/kovo/venv/bin/python -m uvicorn src.gateway.main:app --host 0.0.0.0 --port 8080 --log-level info
Restart=always
RestartSec=5
TimeoutStopSec=30
StandardOutput=append:/opt/kovo/logs/gateway.log
StandardError=append:/opt/kovo/logs/gateway.log
[Install]
WantedBy=multi-user.target
